# marketmaker/files/profile.d/komodo-ac-cli.sh
# aliases for komodo asset chains, cli
alias revs-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/REVS/REVS.conf"
alias supernet-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/SUPERNET/SUPERNET.conf"
alias dex-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/DEX/DEX.conf"
alias pangea-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/PANGEA/PANGEA.conf"
alias jumblr-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/JUMBLR/JUMBLR.conf"
alias bet-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/BET/BET.conf"
alias crypto-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/CRYPTO/CRYPTO.conf"
alias hodl-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/HODL/HODL.conf"
alias mshark-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/MSHARK/MSHARK.conf"
alias bots-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/BOTS/BOTS.conf"
alias mgw-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/MGW/MGW.conf"
alias coqui-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/COQUI/COQUI.conf"
alias wlc-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/WLC/WLC.conf"
alias kv-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/KV/KV.conf"
alias ceal-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/CEAL/CEAL.conf"
alias mesh-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/MESH/MESH.conf"
alias mnz-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/MNZ/MNZ.conf"
alias axo-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/AXO/AXO.conf"
alias btch-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/BTCH/BTCH.conf"
alias pizza-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/PIZA/PIZZA.conf"
alias beer-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/BEER/BEER.conf"
alias ninja-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/NINJA/NINJA.conf"
alias oot-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/OOT/OOT.conf"
alias bntn-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/BNTN/BNTN.conf"
alias chain-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/CHAIN/CHAIN.conf"
alias prlpay-cli="/home/barterdex/bin/komodo-cli -conf=/home/barterdex/.komodo/PRLPAY/PRLPAY.conf"
